# cassie-system/cassie-watchdog.timer
[Unit]
Description=Check Cassie pipeline health every 2 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=120
AccuracySec=10

[Install]
WantedBy=timers.target
